# eos-reboot-required.timer
[Unit]
Description=Notify about reboot on special conditions.

[Timer]
OnCalendar=minutely
Persistent=true
 
[Install]
WantedBy=timers.target
